# Docs/ollama_proxy_exporter.service
[Unit]
Description=Ollama Proxy Exporter
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=raja
Group=raja
WorkingDirectory=/home/raja/ollama-monitoring
ExecStart=/usr/bin/python3 /home/raja/ollama-monitoring/ollama_proxy_exporter.py --listen 0.0.0.0 --port 11435 --backend-host 127.0.0.1 --backend-port 11434 --api-base http://127.0.0.1:11434
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target